# systemd/nsddns.timer
[Unit]
Description=Run nsddns every hour to keep DNS record up to date

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
